# docker/supervisor.conf
[supervisord]
nodaemon = true

[program:nginx]
user = root
command = service nginx restart
